# crontab
# 매일 저녁 9시에 main.py 실행
00 23 * * * /usr/local/bin/python /app/main.py >> /var/log/cron.log 2>&1

